# scripts/supervisor.agent-gateway.conf
# Supervisor config for the FastMCP Agent Gateway.
# Copy to /etc/supervisor/conf.d/ on the host where your CLI agents (claude,
# codex, opencode, hermes) live. The gateway is OFF by default in dev containers
# — only start it where you have the CLIs installed.
#
#   sudo cp /app/scripts/supervisor.agent-gateway.conf /etc/supervisor/conf.d/
#   sudo supervisorctl reread && sudo supervisorctl update
#   sudo supervisorctl start agent-gateway
#
# Set environment via a real .env file at /opt/baseline/agent-gateway/.env and
# point EnvironmentFile= at it (systemd) or use the `environment=` line below.

[program:agent-gateway]
command=/opt/agent-gateway-venv/bin/agent-gateway --host 127.0.0.1 --port 8765
directory=/app/services/agent-gateway
autostart=false
autorestart=true
startsecs=5
stderr_logfile=/var/log/supervisor/agent-gateway.err.log
stdout_logfile=/var/log/supervisor/agent-gateway.out.log
stopsignal=TERM
stopwaitsecs=15
stopasgroup=true
killasgroup=true
environment=GATEWAY_DATA_DIR="/var/lib/agent-gateway",ENABLED_AGENTS="claude,codex,opencode,hermes",GATEWAY_WORKSPACE_ID="1",MC_URL="",MC_API_KEY=""
